# symlinks/global/.config/autostart/stjerm.desktop
[Desktop Entry]
Type=Application
Exec=stjermstart
Hidden=false
NoDisplay=false
X-GNOME-Autostart-enabled=true
Name=stjerm
