# symlinks/.config/autostart/devilspie.desktop
[Desktop Entry]
Type=Application
Name=Devilspie
X-GNOME-Autostart-enabled=true
Exec=timeout 6 devilspie

