# symlinks/adg/.config/autostart/dbus.desktop
[Desktop Entry]
Type=Application
Name=Export dbus
Exec=/home/adg/.bin/export_dbus.sh
NoDisplay=false
X-GNOME-Autostart-enabled=true
